# crontab
# Run once every hour
@hourly node /opt/beeDocker/cronUserDel.js
